# ai.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.gateway</string>
    <key>Comment</key>
    <string>Telegram to Codex Gateway</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>1</integer>
    <key>Umask</key>
    <integer>63</integer>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>plist_path="$(/bin/launchctl print "gui/$(/usr/bin/id -u)/ai.gateway" 2>/dev/null | /usr/bin/sed -n 's/^[[:space:]]*path = //p' | /usr/bin/head -n 1)"; [[ -n "$plist_path" ]] || plist_path="$HOME/Library/LaunchAgents/ai.gateway.plist"; exec "${plist_path:A:h}/launch"</string>
    </array>
  </dict>
</plist>
